# haproxy_exporter.service
[Service]
ExecStart=/usr/bin/haproxy_exporter @/var/lib/haproxy_exporter/args
DynamicUser=yes
StateDirectory=haproxy_exporter
